# data/ayatana-indicator-application.conf.in
description "Ayatana Indicator Application Service"

start on ayatana-indicator-services-start
stop on desktop-end or ayatana-indicator-services-end

emits ayatana-appindicators-start ayatana-appindicators-end

respawn
respawn limit 2 10

pre-start script
	# NOTE: Limiting only to Unity 7 right now as it's still using
	# dbusmenu.  That can be lifted after it is ported to GMenu

	if [ "x$DESKTOP_SESSION" = "xubuntu-touch" ] ; then
		stop; exit 0
	fi
end script

exec $SNAP@pkglibexecdir@/ayatana-indicator-application-service

pre-stop exec initctl emit ayatana-appindicators-end
