# data/ayatana-indicator-application.service.in
[Unit]
Description=Ayatana Indicator Application Service
PartOf=graphical-session.target
PartOf=ayatana-indicators.target

[Service]
ExecStart=@pkglibexecdir@/ayatana-indicator-application-service
Restart=on-failure

[Install]
WantedBy=ayatana-indicators.target
