# data/ayatana-indicator-application.service.in
[Unit]
Description=Ayatana Indicator Application Service
PartOf=graphical-session.target
After=ayatana-indicators-pre.target

[Service]
ExecStart=@pkglibexecdir@/ayatana-indicator-application-service
Restart=on-failure
