# data/ayatana-indicator-application.service.in
[Unit]
Description=Ayatana Indicator Application Service
PartOf=graphical-session.target
PartOf=ayatana-indicators.target

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/@CMAKE_INSTALL_LIBEXECDIR@/ayatana-indicator-application/ayatana-indicator-application-service
Restart=on-failure

[Install]
WantedBy=ayatana-indicators.target
